# unittest/hadoopDocker/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:dfs]
command=/bin/bash -c /opt/hadoop/sbin/start-dfs.sh

[program:yarn]
command=/bin/bash -c /opt/hadoop/sbin/start-yarn.sh

[program:zookeeper]
command=/bin/bash -c /opt/zookeeper/bin/zkServer.sh

[program:hbase]
command=/bin/bash -c /opt/hbase/bin/start-hbase.sh

[program:spark]
command=/bin/bash -c /opt/spark/sbin/start-all.sh
